# auto-avahi.service
[Unit]
Description=Auto-Avahi: TLS cert + Traefik config automation for Docker (with optional mDNS)
Documentation=https://github.com/Tristan-Wilson/auto-avahi
After=docker.service avahi-daemon.service
Requires=docker.service
Wants=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-avahi
Restart=always
RestartSec=10

# Configuration — update these for your environment
# Run: auto-avahi --print-env to see current values

# Required when MDNS_ENABLE=true; ignored otherwise.
Environment="AVAHI_HOST_IP="

# Comma-separated list of allowed hostname suffixes (without leading dot).
# Hostnames extracted from Traefik labels must end with one of these.
# Default behavior preserves the original .local-only mode.
Environment="HOSTNAME_SUFFIXES=local"

# Set to "false" to skip mDNS publishing entirely (useful when DNS is
# handled by another resolver). When false, AVAHI_HOST_IP is not required
# and no avahi-publish processes are spawned.
Environment="MDNS_ENABLE=true"

Environment="CERT_DIR="
Environment="CERT_DIR_CONTAINER=/certs"
Environment="TRAEFIK_CONFIG_DIR="

# Optional: mkcert CA root directory
# Set this if running as root but browsers trust a different user's CA.
# Find the path: mkcert -CAROOT  (run as that user)
#Environment="CAROOT=/home/youruser/.local/share/mkcert"

# Security
User=root
Group=docker

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto-avahi

[Install]
WantedBy=multi-user.target
